# deploy/systemd/reap-capture@.service
[Unit]
Description=Reap credential-free OKX public capture (%i)
Documentation=file:/usr/local/share/doc/reap/operations.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=2

[Service]
Type=simple
User=reap
Group=reap
WorkingDirectory=/var/lib/reap/capture/%i
StateDirectory=reap/connectivity
StateDirectoryMode=0750
Environment=RUST_BACKTRACE=1
EnvironmentFile=/etc/reap/capture/%i.env
ExecStart=/usr/local/bin/reap capture --config /etc/reap/capture/okx-btc-public.toml --output /var/lib/reap/capture/%i/run-report.json --raw-path /var/lib/reap/capture/%i/raw.jsonl --duration-secs ${REAP_CAPTURE_DURATION_SECS} --require-clean-capture
Restart=no
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=45s
UMask=0077
LimitNOFILE=65536
TasksMax=512
NoNewPrivileges=true
CapabilityBoundingSet=
ProtectProc=invisible
PrivateIPC=true
RemoveIPC=true
KeyringMode=private
RestrictRealtime=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/reap/capture/%i /var/lib/reap/connectivity
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reap-capture-%i

[Install]
WantedBy=multi-user.target
